# eggs/systemd/restic-backup.timer
[Unit]
Description=Daily backup timer

[Timer]
OnCalendar=Daily
Persistent=True

[Install]
WantedBy=timers.target
